# bash/cleanup-journal/cleanup-journal.timer
[Unit]
Description=Timer to trim systemd journal

[Timer]
OnBootSec=1hour
OnUnitActiveSec=1hour

[Install]
WantedBy=timers.target
